# supervisord.conf
[supervisord]
user=root
nodaemon=true
pidfile=/var/run/supervisord.pid
logfile=/var/log/supervisor/supervisord.log

[program:serve]
command=%(ENV_SUPERVISOR_SERVE_COMMAND)s
user=%(ENV_SUPERVISOR_SERVE_USER)s
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
